# etc/systemd/system/mecamera-watchdog.service
[Unit]
Description=ME Camera health watchdog
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/mecamera-watchdog.sh auto http://127.0.0.1:8080/api/health
