# etc/systemd/system/mecamera.service
[Unit]
Description=ME Camera Service
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/ME_CAM-DEV
ExecStart=/usr/bin/python3 /home/pi/ME_CAM-DEV/main.py
Environment="PYTHONPATH=/home/pi/.local/lib/python3.9/site-packages"
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
